# tire-agent-web.service
[Unit]
Description=Tire Sales Agent – Frontend (Node.js/Express)
After=network.target tire-agent.service
Requires=tire-agent.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tire_sales_agent
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tire-agent-web

[Install]
WantedBy=multi-user.target
